# .config/systemd/user/wiki-backup.service
[Unit]
Description=Backup personal WIKI

[Service]
Type=exec
ExecStart=%h/projects/wiki/upload.sh

[Install]
WantedBy=default.target
